# MGNREGA_District_Pulse_Karnataka_full/infra/systemd/etl.timer
[Unit]
Description=Run MGNREGA ETL daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
